# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:laravel-setup]
command=bash -c "php artisan migrate --force && php artisan storage:link && php artisan config:cache && php artisan queue:restart"
directory=/app
autostart=true
autorestart=false
startsecs=0
startretries=1
user=root
stdout_logfile=/tmp/setup.log
stderr_logfile=/tmp/setup_error.log
priority=1

[program:laravel-worker]
command=php artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --memory=512 --timeout=300
directory=/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/tmp/worker.log
stderr_logfile=/tmp/worker_error.log
stopwaitsecs=10
priority=2

[program:laravel-server]
command=php artisan serve --host=0.0.0.0 --port=%(ENV_PORT)s
directory=/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
stdout_logfile=/tmp/server.log
stderr_logfile=/tmp/server_error.log
priority=3